# medusa-0.5.1/index-configuration/medusa.cron.in
#!/bin/sh

(
    sleep 600

    indexing_is_enabled=`@prefix@/sbin/medusa-test-indexing-status`

    if [ "$indexing_is_enabled" = "Indexing is enabled." ]
    then
	killall medusa-indexd
	killall medusa-searchd
	@prefix@/bin/medusa-searchd
	@prefix@/bin/medusa-indexd
    fi
) &

# Later we'll also remove index files here, if medusa is not enabled

